# Ebuilds/forge-pvp/files/listen-forge-pvp.socket
[Unit]
Description=Listen the choosing port of forge-pvp

[Socket]
ListenStream=11026

[Install]
WantedBy=multi-user.target
