# Ebuilds/forge-pvp/files/stop-forge-pvp.timer
[Timer]
OnUnitActiveSec=60s
Unit=stop-forge-pvp.service

[Install]
WantedBy=
