# Ebuilds/minecraft-server-script/files/stop-fabric-survival.timer
[Timer]
OnUnitActiveSec=60s
Unit=stop-fabric-survival.service

[Install]
WantedBy=
