# Ebuilds/fabric-survival/files/listen-fabric-survival.socket
[Unit]
Description=Listen the choosing port of fabric-survival

[Socket]
ListenStream=11024

[Install]
WantedBy=multi-user.target
