# Ebuilds/vanilla-survival/files/listen-vanilla-survival.socket
[Unit]
Description=Listen the choosing port of vanilla-survival

[Socket]
ListenStream=11028

[Install]
WantedBy=multi-user.target
