# Ebuilds/vanilla-survival/files/stop-vanilla-survival.timer
[Timer]
OnUnitActiveSec=60s
Unit=stop-vanilla-survival.service

[Install]
WantedBy=
